# plejd/rootfs/etc/services.d/plejd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Plejd (TypeScript)
# Runs the Plejd service
# ==============================================================================

bashio::log.info 'Starting the Plejd TypeScript service...'

# Ensure proper permissions
chmod +x /usr/bin/plejd.sh

# Execute the service
exec /usr/bin/plejd.sh