# plejd/rootfs/etc/services.d/plejd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Plejd
# Runs the Plejd service
# ==============================================================================

bashio::log.info 'Starting the Plejd service...'

# Change working directory
cd /plejd || bashio::exit.nok 'Unable to change working directory'

# Run the Plejd service
bashio::log.info 'Verified permissions on startup script'
chmod +x /usr/bin/plejd.sh

bashio::log.info 'Executing startup script'
exec /usr/bin/plejd.sh
